# dicom-gateway/production/systemd/rai-orthanc-import.service
[Unit]
Description=RAI PACS Orthanc pending study import
Wants=network-online.target
After=network-online.target docker.service

[Service]
Type=oneshot
User=support
Group=support
WorkingDirectory=/opt/rai-pacs
Environment=RAI_PACS_GATEWAY_ENV_FILE=/opt/rai-pacs/dicom-gateway/production/.env
Environment=RAI_PACS_ORTHANC_IMPORT_STATE_FILE=/opt/rai-pacs/dicom-gateway/production/.orthanc-import-state.json
Environment=RAI_PACS_ORTHANC_EVENT_STATE_FILE=/opt/rai-pacs/dicom-gateway/production/.orthanc-event-state.json
Environment=RAI_PACS_ORTHANC_LOG_STATE_FILE=/opt/rai-pacs/dicom-gateway/production/.orthanc-log-state.json
ExecStartPre=/usr/bin/npm run sync:orthanc-logs
ExecStartPre=/usr/bin/npm run sync:orthanc-events
ExecStart=/usr/bin/npm run import:orthanc-pending
