# dicom-gateway/production/systemd/rai-orthanc-import.timer
[Unit]
Description=Run RAI PACS Orthanc import watcher every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
